# app/rootfs/etc/supervisor.d/shopware-scheduled-task.ini
[program:shopware-scheduled-task]
command=/usr/local/bin/php -d "memory_limit=2G" /var/www/html/bin/console scheduled-task:run --time-limit=300 --memory-limit=512M
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
autostart=true
autorestart=true
environment=HOME="/var/www/html",USER="www-data"
